# systemd/iot.service
[Unit]
Description=IOT-Service

[Service]
Type=simple
ExecStart=/opt/iotprojects/bin/webservices.sh

[Install]
WantedBy=multi-user.target